# desktop/.config/systemd/user/hyprpaper.service
[Unit]
Description=hyprpaper

[Service]
Type=simple
ExecStart=/usr/bin/hyprpaper

[Install]
WantedBy=default.target
