# desktop/.config/systemd/user/lxpolkit.service
[Unit]
Description=lxpolkit

[Service]
Type=simple
ExecStart=/usr/bin/lxpolkit

[Install]
WantedBy=default.target
